# setup/etc/cron.daily/CLEARSTASH
#!/bin/bash
python /usr/scripts/clearStash/runJob.py
